# covid-terminator/covid-terminator-0.1/debian/covid-terminator.target
[Unit]
Description=Covid Terminator Target
Documentation=man:systemd:special(7)
Requires=covid-post-terminate.service
After=covid-post-terminate.service
AllowIsolate=no
